# templates/default/systemd.service.erb
[Unit]
Description=Confd
After=network.target

[Service]
User=root
ExecStart=/usr/local/bin/confd -config-file <%= node['confd']['confdir'] %>/confd.toml

[Install]
WantedBy=multi-user.target
